# deploy/veriloga-mcp.service
[Unit]
Description=VerilogA MCP Server (veriloga-help)
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=mcp
Group=mcp
WorkingDirectory=/opt/mcp/veriloga-help
ExecStart=/opt/mcp/veriloga-help/.venv/bin/python server/main.py --host 0.0.0.0 --port 8096
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=veriloga-mcp

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
